# docker/drupal/rootfs/etc/services.d/php/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/php-fpm -y /etc/php/${PHP_VERSION}/fpm/php-fpm.conf
